# sysconf/etc/profile.d/sudowarn.sh
alias sudo=/bin/sudo
alias pkexec=/bin/sudo